# backup.service
[Unit]
Description=Backup with rclone

[Service]
Type=simple
ExecStart=%h/bin/backup.sh

[Install]
WantedBy=default.target
